# dev-tools/src/main/vagrant/develop/etc/profile.d/kapua_env.sh
export ACTIVEMQ_OPTS="-Dcommons.db.schema.update=true"
export CATALINA_OPTS="-Dcommons.db.schema.update=true"
